# steam.container
[Unit]
Description=Steam Headless (gamescope) in Podman

[Container]
Image=localhost/arch-steam-headless
UserNS=keep-id
User=%U
Device=/dev/dri
Network=host
IPC=host
ShmSize=2g
Environment=HOME=/home/steamuser
# Tweak virtual display and FPS here:
Environment=GAMESCOPE_WIDTH=1920
Environment=GAMESCOPE_HEIGHT=1080
Environment=GAMESCOPE_FPS=60
# Persisted volumes:
Volume=%h/docker-steam/home:/home/steamuser:Z
Volume=%h/docker-steam/.steam:/home/steamuser/.steam:Z
Volume=%h/docker-steam/.local/share/Steam:/home/steamuser/.local/share/Steam:Z
# Default command already starts headless:
Exec=/usr/local/bin/start-headless.sh
Name=steam-headless

[Install]
WantedBy=default.target
